# .config/systemd/user/noir-sweep.service
[Unit]
Description=noir sunday sweep — weekly agent pass over the vault

[Service]
Type=oneshot
ExecStart=%h/.bin/sweep.sh
