# .config/systemd/user/noir-brief.service
[Unit]
Description=noir morning brief — seed today's journal

[Service]
Type=oneshot
ExecStart=%h/.bin/brief.sh morning
